# loki-storage-servers.target
[Unit]
Description=loki storage servers
