# oxen-storage-servers.target
[Unit]
Description=Oxen storage servers
