# oxen-nodes.target
[Unit]
Description=loki nodes

[Install]
Alias=loki-nodes.target
